# systemd/librescoot-alarm.service
[Unit]
Description=LibreScoot Alarm Service
After=network.target redis.service bmx-service.service
Wants=redis.service bmx-service.service

[Service]
Type=simple
ExecStart=/usr/bin/alarm-service \
    --redis=localhost:6379 \
    --log-level=info
Restart=always
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target