# systemd/librescoot-alarm.service
[Unit]
Description=LibreScoot Alarm Service
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
ExecStart=/usr/bin/alarm-service \
    --i2c-bus=/dev/i2c-3 \
    --redis=localhost:6379 \
    --log-level=info \
    --alarm-duration=10
Restart=always
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target